# playbooks/templates/pidisplay.service.j2
[Unit]
Description=Pi Display Service
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python3 {{ pidi_script_folder }}/pi4display.py

[Install]
WantedBy=multi-user.target
